# contrib/posix/templates/systemd.timer.in
{{/*
# Template for generating systemd timers
# Used to schedule commands when "scheduler=systemd" (default)
*/}}

[Unit]
Description={{ .TimerDescription }}

[Timer]
{{ range .OnCalendar -}}
OnCalendar={{ . }}
{{ end -}}
Unit={{ .SystemdProfile }}
Persistent=true

[Install]
WantedBy=timers.target
