# schtasks/examples/every_hour_on_Mondays_in_january.xml
<?xml version="1.0" encoding="UTF-16"?>
<Task version="1.2" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Date>2025-02-06T22:46:30</Date>
    <Author>resticprofile</Author>
    <Description>every hour on Mondays in january</Description>
    <URI>\resticprofile backup\test 17</URI>
  </RegistrationInfo>
  <Principals>
    <Principal id="Author">
      <UserId>S-1-5-21-3622340250-3004114779-3108145013-1002</UserId>
      <LogonType>InteractiveToken</LogonType>
    </Principal>
  </Principals>
  <Settings>
    <DisallowStartIfOnBatteries>true</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>true</StopIfGoingOnBatteries>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <IdleSettings>
      <Duration>PT10M</Duration>
      <WaitTimeout>PT1H</WaitTimeout>
      <StopOnIdleEnd>true</StopOnIdleEnd>
      <RestartOnIdle>false</RestartOnIdle>
    </IdleSettings>
  </Settings>
  <Triggers>
    <CalendarTrigger>
      <StartBoundary>2026-01-05T00:04:00</StartBoundary>
      <ExecutionTimeLimit>PT0S</ExecutionTimeLimit>
      <ScheduleByMonthDayOfWeek>
        <Months>
          <January />
        </Months>
        <Weeks>
          <Week>1</Week>
          <Week>2</Week>
          <Week>3</Week>
          <Week>4</Week>
          <Week>Last</Week>
        </Weeks>
        <DaysOfWeek>
          <Monday />
        </DaysOfWeek>
      </ScheduleByMonthDayOfWeek>
    </CalendarTrigger>
    <CalendarTrigger>
      <StartBoundary>2026-01-05T01:04:00</StartBoundary>
      <ExecutionTimeLimit>PT0S</ExecutionTimeLimit>
      <ScheduleByMonthDayOfWeek>
        <Months>
          <January />
        </Months>
        <Weeks>
          <Week>1</Week>
          <Week>2</Week>
          <Week>3</Week>
          <Week>4</Week>
          <Week>Last</Week>
        </Weeks>
        <DaysOfWeek>
          <Monday />
        </DaysOfWeek>
      </ScheduleByMonthDayOfWeek>
    </CalendarTrigger>
    <CalendarTrigger>
      <StartBoundary>2026-01-05T02:04:00</StartBoundary>
      <ExecutionTimeLimit>PT0S</ExecutionTimeLimit>
      <ScheduleByMonthDayOfWeek>
        <Months>
          <January />
        </Months>
        <Weeks>
          <Week>1</Week>
          <Week>2</Week>
          <Week>3</Week>
          <Week>4</Week>
          <Week>Last</Week>
        </Weeks>
        <DaysOfWeek>
          <Monday />
        </DaysOfWeek>
      </ScheduleByMonthDayOfWeek>
    </CalendarTrigger>
    <CalendarTrigger>
      <StartBoundary>2026-01-05T03:04:00</StartBoundary>
      <ExecutionTimeLimit>PT0S</ExecutionTimeLimit>
      <ScheduleByMonthDayOfWeek>
        <Months>
          <January />
        </Months>
        <Weeks>
          <Week>1</Week>
          <Week>2</Week>
          <Week>3</Week>
          <Week>4</Week>
          <Week>Last</Week>
        </Weeks>
        <DaysOfWeek>
          <Monday />
        </DaysOfWeek>
      </ScheduleByMonthDayOfWeek>
    </CalendarTrigger>
    <CalendarTrigger>
      <StartBoundary>2026-01-05T04:04:00</StartBoundary>
      <ExecutionTimeLimit>PT0S</ExecutionTimeLimit>
      <ScheduleByMonthDayOfWeek>
        <Months>
          <January />
        </Months>
        <Weeks>
          <Week>1</Week>
          <Week>2</Week>
          <Week>3</Week>
          <Week>4</Week>
          <Week>Last</Week>
        </Weeks>
        <DaysOfWeek>
          <Monday />
        </DaysOfWeek>
      </ScheduleByMonthDayOfWeek>
    </CalendarTrigger>
    <CalendarTrigger>
      <StartBoundary>2026-01-05T05:04:00</StartBoundary>
      <ExecutionTimeLimit>PT0S</ExecutionTimeLimit>
      <ScheduleByMonthDayOfWeek>
        <Months>
          <January />
        </Months>
        <Weeks>
          <Week>1</Week>
          <Week>2</Week>
          <Week>3</Week>
          <Week>4</Week>
          <Week>Last</Week>
        </Weeks>
        <DaysOfWeek>
          <Monday />
        </DaysOfWeek>
      </ScheduleByMonthDayOfWeek>
    </CalendarTrigger>
    <CalendarTrigger>
      <StartBoundary>2026-01-05T06:04:00</StartBoundary>
      <ExecutionTimeLimit>PT0S</ExecutionTimeLimit>
      <ScheduleByMonthDayOfWeek>
        <Months>
          <January />
        </Months>
        <Weeks>
          <Week>1</Week>
          <Week>2</Week>
          <Week>3</Week>
          <Week>4</Week>
          <Week>Last</Week>
        </Weeks>
        <DaysOfWeek>
          <Monday />
        </DaysOfWeek>
      </ScheduleByMonthDayOfWeek>
    </CalendarTrigger>
    <CalendarTrigger>
      <StartBoundary>2026-01-05T07:04:00</StartBoundary>
      <ExecutionTimeLimit>PT0S</ExecutionTimeLimit>
      <ScheduleByMonthDayOfWeek>
        <Months>
          <January />
        </Months>
        <Weeks>
          <Week>1</Week>
          <Week>2</Week>
          <Week>3</Week>
          <Week>4</Week>
          <Week>Last</Week>
        </Weeks>
        <DaysOfWeek>
          <Monday />
        </DaysOfWeek>
      </ScheduleByMonthDayOfWeek>
    </CalendarTrigger>
    <CalendarTrigger>
      <StartBoundary>2026-01-05T08:04:00</StartBoundary>
      <ExecutionTimeLimit>PT0S</ExecutionTimeLimit>
      <ScheduleByMonthDayOfWeek>
        <Months>
          <January />
        </Months>
        <Weeks>
          <Week>1</Week>
          <Week>2</Week>
          <Week>3</Week>
          <Week>4</Week>
          <Week>Last</Week>
        </Weeks>
        <DaysOfWeek>
          <Monday />
        </DaysOfWeek>
      </ScheduleByMonthDayOfWeek>
    </CalendarTrigger>
    <CalendarTrigger>
      <StartBoundary>2026-01-05T09:04:00</StartBoundary>
      <ExecutionTimeLimit>PT0S</ExecutionTimeLimit>
      <ScheduleByMonthDayOfWeek>
        <Months>
          <January />
        </Months>
        <Weeks>
          <Week>1</Week>
          <Week>2</Week>
          <Week>3</Week>
          <Week>4</Week>
          <Week>Last</Week>
        </Weeks>
        <DaysOfWeek>
          <Monday />
        </DaysOfWeek>
      </ScheduleByMonthDayOfWeek>
    </CalendarTrigger>
    <CalendarTrigger>
      <StartBoundary>2026-01-05T10:04:00</StartBoundary>
      <ExecutionTimeLimit>PT0S</ExecutionTimeLimit>
      <ScheduleByMonthDayOfWeek>
        <Months>
          <January />
        </Months>
        <Weeks>
          <Week>1</Week>
          <Week>2</Week>
          <Week>3</Week>
          <Week>4</Week>
          <Week>Last</Week>
        </Weeks>
        <DaysOfWeek>
          <Monday />
        </DaysOfWeek>
      </ScheduleByMonthDayOfWeek>
    </CalendarTrigger>
    <CalendarTrigger>
      <StartBoundary>2026-01-05T11:04:00</StartBoundary>
      <ExecutionTimeLimit>PT0S</ExecutionTimeLimit>
      <ScheduleByMonthDayOfWeek>
        <Months>
          <January />
        </Months>
        <Weeks>
          <Week>1</Week>
          <Week>2</Week>
          <Week>3</Week>
          <Week>4</Week>
          <Week>Last</Week>
        </Weeks>
        <DaysOfWeek>
          <Monday />
        </DaysOfWeek>
      </ScheduleByMonthDayOfWeek>
    </CalendarTrigger>
    <CalendarTrigger>
      <StartBoundary>2026-01-05T12:04:00</StartBoundary>
      <ExecutionTimeLimit>PT0S</ExecutionTimeLimit>
      <ScheduleByMonthDayOfWeek>
        <Months>
          <January />
        </Months>
        <Weeks>
          <Week>1</Week>
          <Week>2</Week>
          <Week>3</Week>
          <Week>4</Week>
          <Week>Last</Week>
        </Weeks>
        <DaysOfWeek>
          <Monday />
        </DaysOfWeek>
      </ScheduleByMonthDayOfWeek>
    </CalendarTrigger>
    <CalendarTrigger>
      <StartBoundary>2026-01-05T13:04:00</StartBoundary>
      <ExecutionTimeLimit>PT0S</ExecutionTimeLimit>
      <ScheduleByMonthDayOfWeek>
        <Months>
          <January />
        </Months>
        <Weeks>
          <Week>1</Week>
          <Week>2</Week>
          <Week>3</Week>
          <Week>4</Week>
          <Week>Last</Week>
        </Weeks>
        <DaysOfWeek>
          <Monday />
        </DaysOfWeek>
      </ScheduleByMonthDayOfWeek>
    </CalendarTrigger>
    <CalendarTrigger>
      <StartBoundary>2026-01-05T14:04:00</StartBoundary>
      <ExecutionTimeLimit>PT0S</ExecutionTimeLimit>
      <ScheduleByMonthDayOfWeek>
        <Months>
          <January />
        </Months>
        <Weeks>
          <Week>1</Week>
          <Week>2</Week>
          <Week>3</Week>
          <Week>4</Week>
          <Week>Last</Week>
        </Weeks>
        <DaysOfWeek>
          <Monday />
        </DaysOfWeek>
      </ScheduleByMonthDayOfWeek>
    </CalendarTrigger>
    <CalendarTrigger>
      <StartBoundary>2026-01-05T15:04:00</StartBoundary>
      <ExecutionTimeLimit>PT0S</ExecutionTimeLimit>
      <ScheduleByMonthDayOfWeek>
        <Months>
          <January />
        </Months>
        <Weeks>
          <Week>1</Week>
          <Week>2</Week>
          <Week>3</Week>
          <Week>4</Week>
          <Week>Last</Week>
        </Weeks>
        <DaysOfWeek>
          <Monday />
        </DaysOfWeek>
      </ScheduleByMonthDayOfWeek>
    </CalendarTrigger>
    <CalendarTrigger>
      <StartBoundary>2026-01-05T16:04:00</StartBoundary>
      <ExecutionTimeLimit>PT0S</ExecutionTimeLimit>
      <ScheduleByMonthDayOfWeek>
        <Months>
          <January />
        </Months>
        <Weeks>
          <Week>1</Week>
          <Week>2</Week>
          <Week>3</Week>
          <Week>4</Week>
          <Week>Last</Week>
        </Weeks>
        <DaysOfWeek>
          <Monday />
        </DaysOfWeek>
      </ScheduleByMonthDayOfWeek>
    </CalendarTrigger>
    <CalendarTrigger>
      <StartBoundary>2026-01-05T17:04:00</StartBoundary>
      <ExecutionTimeLimit>PT0S</ExecutionTimeLimit>
      <ScheduleByMonthDayOfWeek>
        <Months>
          <January />
        </Months>
        <Weeks>
          <Week>1</Week>
          <Week>2</Week>
          <Week>3</Week>
          <Week>4</Week>
          <Week>Last</Week>
        </Weeks>
        <DaysOfWeek>
          <Monday />
        </DaysOfWeek>
      </ScheduleByMonthDayOfWeek>
    </CalendarTrigger>
    <CalendarTrigger>
      <StartBoundary>2026-01-05T18:04:00</StartBoundary>
      <ExecutionTimeLimit>PT0S</ExecutionTimeLimit>
      <ScheduleByMonthDayOfWeek>
        <Months>
          <January />
        </Months>
        <Weeks>
          <Week>1</Week>
          <Week>2</Week>
          <Week>3</Week>
          <Week>4</Week>
          <Week>Last</Week>
        </Weeks>
        <DaysOfWeek>
          <Monday />
        </DaysOfWeek>
      </ScheduleByMonthDayOfWeek>
    </CalendarTrigger>
    <CalendarTrigger>
      <StartBoundary>2026-01-05T19:04:00</StartBoundary>
      <ExecutionTimeLimit>PT0S</ExecutionTimeLimit>
      <ScheduleByMonthDayOfWeek>
        <Months>
          <January />
        </Months>
        <Weeks>
          <Week>1</Week>
          <Week>2</Week>
          <Week>3</Week>
          <Week>4</Week>
          <Week>Last</Week>
        </Weeks>
        <DaysOfWeek>
          <Monday />
        </DaysOfWeek>
      </ScheduleByMonthDayOfWeek>
    </CalendarTrigger>
    <CalendarTrigger>
      <StartBoundary>2026-01-05T20:04:00</StartBoundary>
      <ExecutionTimeLimit>PT0S</ExecutionTimeLimit>
      <ScheduleByMonthDayOfWeek>
        <Months>
          <January />
        </Months>
        <Weeks>
          <Week>1</Week>
          <Week>2</Week>
          <Week>3</Week>
          <Week>4</Week>
          <Week>Last</Week>
        </Weeks>
        <DaysOfWeek>
          <Monday />
        </DaysOfWeek>
      </ScheduleByMonthDayOfWeek>
    </CalendarTrigger>
    <CalendarTrigger>
      <StartBoundary>2026-01-05T21:04:00</StartBoundary>
      <ExecutionTimeLimit>PT0S</ExecutionTimeLimit>
      <ScheduleByMonthDayOfWeek>
        <Months>
          <January />
        </Months>
        <Weeks>
          <Week>1</Week>
          <Week>2</Week>
          <Week>3</Week>
          <Week>4</Week>
          <Week>Last</Week>
        </Weeks>
        <DaysOfWeek>
          <Monday />
        </DaysOfWeek>
      </ScheduleByMonthDayOfWeek>
    </CalendarTrigger>
    <CalendarTrigger>
      <StartBoundary>2026-01-05T22:04:00</StartBoundary>
      <ExecutionTimeLimit>PT0S</ExecutionTimeLimit>
      <ScheduleByMonthDayOfWeek>
        <Months>
          <January />
        </Months>
        <Weeks>
          <Week>1</Week>
          <Week>2</Week>
          <Week>3</Week>
          <Week>4</Week>
          <Week>Last</Week>
        </Weeks>
        <DaysOfWeek>
          <Monday />
        </DaysOfWeek>
      </ScheduleByMonthDayOfWeek>
    </CalendarTrigger>
    <CalendarTrigger>
      <StartBoundary>2026-01-05T23:04:00</StartBoundary>
      <ExecutionTimeLimit>PT0S</ExecutionTimeLimit>
      <ScheduleByMonthDayOfWeek>
        <Months>
          <January />
        </Months>
        <Weeks>
          <Week>1</Week>
          <Week>2</Week>
          <Week>3</Week>
          <Week>4</Week>
          <Week>Last</Week>
        </Weeks>
        <DaysOfWeek>
          <Monday />
        </DaysOfWeek>
      </ScheduleByMonthDayOfWeek>
    </CalendarTrigger>
  </Triggers>
  <Actions Context="Author">
    <Exec>
      <Command>echo</Command>
      <Arguments>hello there</Arguments>
      <WorkingDirectory>C:\</WorkingDirectory>
    </Exec>
  </Actions>
</Task>